# systemd/user/ChangeWallPaper.service
[Unit]
Description=My Script Service

[Service]
ExecStart=/home/nicklas/dotfiles/scripts/changeWallPaper.sh
